# example_daemon.service
[Unit]
Description=Health Data Pipeline API
After=network.target postgresql.service
Wants=postgresql.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=ladvien
Group=ladvien
WorkingDirectory=/home/ladvien/self-sensored
Environment=PATH=/home/ladvien/.cache/pypoetry/virtualenvs/self-sensored-tuBp6WIq-py3.13/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/ladvien/self-sensored
ExecStart=/home/ladvien/.cache/pypoetry/virtualenvs/self-sensored-tuBp6WIq-py3.13/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level info
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=health-data-pipeline

[Install]
WantedBy=multi-user.target